# platforms/linux/autostart/ppuc-pinmame.desktop.in
[Desktop Entry]
Type=Application
Version=1.0
Name=PPUC PinMAME
Comment=Start ppuc-pinmame automatically when the GNOME session begins
Terminal=false
Exec=@AUTOSTART_EXEC@
X-GNOME-Autostart-enabled=true
StartupNotify=false
Categories=Game;
